# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ] ; then
	invoke-rc.d dbus reload
fi

#DEBHELPER#

exit 0


